# linux/etc/udev/rules.d/56-g500.rules
SUBSYSTEM=="usb", ATTR{idVendor}=="046d", ATTR{idProduct}=="c068", RUN+="/usr/local/bin/findxinput-g500.sh"
